# deploy/short-bus-jawa.service
[Unit]
Description=Short Bus Jawa Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main
# If the bot doesn't send a watchdog ping (see utils/systemd_watchdog.py)
# within this window, systemd treats it as hung and force-restarts it -
# not just crash detection, but actual hang detection.
WatchdogSec=60
User=jawa
WorkingDirectory=/home/jawa/short-bus-jawa
ExecStart=/home/jawa/short-bus-jawa/venv/bin/python bot.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
